# crontab
@reboot ip link add dummy type dummy
@reboot ip link set dummy up
@reboot ip addr add dev dummy [IPv6 address]
